# pifon/rc.local
#!/bin/sh -e
#
# rc.local file for pifon startup
# copy to /etc and make sure its executable
#

RASPI_USER=chris
RASPI_HOME=/home/chris/raspi
DARKICE=/usr/bin/darkice

# run darkice
if [ -x $DARKICE ]; then
   (cd $RASPI_HOME ; ./tools/tmux_run $RASPI_USER darkice "./tools/looper $DARKICE")
fi

# run pifon_audio
(cd $RASPI_HOME ; ./tools/tmux_run $RASPI_USER audio "./tools/looper ./tools/xmppbot.py ./pifon/pifon_audio") 

exit 0
